# vaani.service
[Unit]
Description=Vaani Speech Assistant
After=network.target

[Service]
Type=simple
ExecStart=/home/$USER/vaani/vaani
WorkingDirectory=/home/$USER/vaani
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target 
